# pi/usbip-allow-sync.path
[Unit]
Description=Watch clients.allow and refresh nftables allow-list
# Replace {{AUTOHUB_ROOT}} with the actual path recorded in /etc/autohub-usbip.conf

[Path]
PathChanged={{AUTOHUB_ROOT}}/config/clients.allow

[Install]
WantedBy=multi-user.target
