# pi/usbip-allow-sync.service
[Unit]
Description=Sync USB/IP client allow-list into nftables
After=network-pre.target
ConditionPathExists=/etc/autohub-usbip.conf

[Service]
Type=oneshot
EnvironmentFile=/etc/autohub-usbip.conf
ExecStart=/bin/bash -c 'exec "$1"' allow "${AUT0HUB_ROOT}/bin/usbip-allow-sync"
